# publisher/sched-tasks/crontab
# Radio-T weekly episode preparation
# Runs every Tuesday at 2:00 PM
0 14 * * 2 cd /srv && publisher prep && publisher deploy